# etc/sv/zookeeper/run
#!/bin/bash
exec 2>&1
exec /usr/share/zookeeper/bin/zkServer.sh start-foreground
