# etc/sv/micro-httpd/run
#!/bin/bash
exec 2>&1
exec /usr/bin/tcpserver 0 80 /usr/sbin/micro-httpd /ion/var/www
