# etc/sv/logstash/run
#!/bin/bash
exec 2>&1
exec /opt/logstash/bin/logstash -f /etc/logstash/hydrometrics.conf -- web
